# dev-embedded/klipper/files/klipper
#!/sbin/openrc-run

command="python klippy"
command_args="${PRINTER_CONFIG} -l ${LOGFILE}"
pidfile="/run/klipper.pid"
